# object-detection-api.service
[Unit]
Description=Object Detection API Service
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/object-detection-api
Environment="PATH=/opt/object-detection-api/venv/bin"
ExecStart=/opt/object-detection-api/venv/bin/gunicorn --config gunicorn_config.py api_server_production:app
Restart=always
RestartSec=10
StandardOutput=append:/opt/object-detection-api/logs/service.log
StandardError=append:/opt/object-detection-api/logs/service.log

[Install]
WantedBy=multi-user.target

